# electrum-server/electrum-server/rw/config/systemd/bitcoind-qrexec.socket
[Unit]
Description=Bitcoind RPC QREXEC Socket

[Socket]
ListenStream=127.0.0.1:8332
Accept=true

[Install]
WantedBy=sockets.target
